# rc.local
#!/bin/sh -e
#
#rc.local
# this file is exetted on system startup
sudo /bin/bash /home/nanocode/RVM/startUP.sh
exit 0
